# config/profile.d/brew.sh
# Linuxbrew PATH for interactive shells
if [ -d /home/linuxbrew/.linuxbrew ]; then
    eval "$(/home/linuxbrew/.linuxbrew/bin/brew shellenv)"
fi
