# config/systemd/bootc-fetch-apply-updates.timer.d/schedule.conf
[Timer]
# Override default schedule: Tuesday 4am Central
# 09:00 UTC = 4am CDT, 10:00 UTC = 4am CST
# The 2hr window covers both DST states
OnCalendar=
OnCalendar=Tue *-*-* 09:00:00
RandomizedDelaySec=1h
