# config/quadlets/ollama.container
[Unit]
Description=Ollama LLM Server (Vulkan)
After=network-online.target magus-provision-data.service
Wants=network-online.target

[Container]
Image=ghcr.io/lazypower/magus-ollama:latest
ContainerName=ollama

# GPU device passthrough
AddDevice=/dev/kfd:/dev/kfd
AddDevice=/dev/dri:/dev/dri

# Group access for GPU devices
GroupAdd=video
GroupAdd=render

# Enable Vulkan GPU acceleration
Environment=OLLAMA_VULKAN=1

# Model storage on secondary NVMe
Volume=/var/data/ollama:/root/.ollama:Z

# API port
PublishPort=11434:11434

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
